# cloudflare-tunnel/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio

config_path=$(bashio::config "CONFIG_PATH")

bashio::log.info "Starting cloudflared..."

exec /opt/cloudflared tunnel --metrics 172.30.32.2:8099 --config "$config_path" run
