# cloudflare-tunnel/rootfs/etc/services.d/cloudflared/run
#!/usr/bin/with-contenv bashio

config_path=$(bashio::config "CONFIG_PATH")

bashio::log.info "Starting cloudflared..."

ls -la /opt/

exec /opt/cloudflared tunnel --metrics 127.0.0.1:8099 --config "$config_path" run
